# docker/server/config/periodic/15min/wpt_cron_call.sh
#!/usr/bin/env sh
/scripts/wpt_cron_call.sh 15min.php
