# docker/server/config/periodic/hourly/wpt_cron_call.sh
#!/usr/bin/env sh
/scripts/wpt_cron_call.sh hourly.php
